# roles/matrix-backup-borg/templates/systemd/matrix-backup-borg.timer.j2
[Unit]
Description=Matrix Borg Backup timer

[Timer]
Unit=matrix-backup-borg.service
OnCalendar={{ matrix_backup_borg_schedule }}
RandomizedDelaySec=2h

[Install]
WantedBy=timers.target
