# roles/matrix-base/templates/systemd/matrix.target
[Unit]
Description=Start all Matrix-related services

[Install]
WantedBy=multi-user.target
