# roles/matrix-coturn/templates/systemd/matrix-coturn-reload.service.j2
[Unit]
Description=Reloads matrix-coturn so that new SSL certificates can kick in

[Service]
Type=oneshot
ExecStart={{ matrix_host_command_systemctl }} reload matrix-coturn.service
