# roles/custom/matrix-synapse-auto-compressor/templates/matrix-synapse-auto-compressor.timer.j2
[Unit]
Description=Synapse State Auto Compressor Timer

[Timer]
Unit=matrix-synapse-auto-compressor.service
OnCalendar={{ matrix_synapse_auto_compressor_calendar }}

[Install]
WantedBy=timers.target
