# roles/matrix-coturn/templates/cron.d/matrix-coturn-ssl-reload.j2
20 4 */5 * * root {{ matrix_host_command_systemctl }} reload matrix-coturn.service
